# netatalk-3.1.11/distrib/initscripts/rc.netbsd
#!/bin/sh
#
# PROVIDE: netatalk
# KEYWORD: shutdown
#
# AFP fileserver for Macintosh clients (netatalk 3.1.11)
#

. /etc/rc.subr

name="netatalk"
rcvar=$name
command="/usr/local/sbin/netatalk"
etcdir="/usr/local/etc"
pidfile="/var/run/${name}.pid"
required_files="$etcdir/afp.conf"

load_rc_config $name
run_rc_command "$1"
